# debian/anerd-server.anerd-udp.upstart
# anerd-udp - a network exchange randomness daemon

description	"a network exchange randomness daemon (udp)"
author		"Dustin Kirkland <dustin.kirkland@gmail.com>"

start on start on runlevel [2345]
stop on start on runlevel [!2345]

respawn
respawn limit 10 5

expect fork

script
	. /etc/default/anerd-server
	if [ -n "$UDP_PORT" ]; then
		exec anerd-udp -p $UDP_PORT
	else
		exec true
	fi
end script
