# debian/anerd-server.anerd-server-tcp.upstart
# anerd-server-tcp - a network exchange randomness daemon

description	"a network exchange randomness daemon (tcp)"
author		"Dustin Kirkland <dustin.kirkland@gmail.com>"

start on start on runlevel [2345]
stop on start on runlevel [!2345]

respawn
respawn limit 10 5

script
	. /etc/default/anerd-server
	if [ -n "$TCP_PORT" ]; then
		exec go run /usr/share/anerd/web/anerd-server-tcp.go $TCP_PORT
	else
		exec true
	fi
end script
