# debian/anerd-server.upstart
# anerd - asynchronous network exchange randomness daemon

description	"asynchronous network exchange randomness daemon"
author		"Dustin Kirkland <kirkland@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

expect fork

# Modify default options here
exec anerd -d /dev/urandom -s 64 -p 26373
