# debian/anerd-client.cron.d
# Run once per minute by default, or tune to your liking here
* * * * *	daemon	run-one anerd-client
