# debian/pollinate.upstart
# pollinate - seed the PRNG using an Entropy-as-a-Service provider

description	"Entropy-as-a-Service client"
author		"Dustin Kirkland <dustin.kirkland@gmail.com>"

setuid daemon

start on (started networking or starting ssh)

script
	. /etc/default/pollinate
	exec pollinate
end script
