# debian/pollinate.upstart
# pollinate - seed the PRNG using an Entropy-as-a-Service provider

description	"Pollinate to seed the pseudo random number generator"
author		"Dustin Kirkland <kirkland@ubuntu.com>"

setuid pollinate

start on (started networking or starting ssh)

exec pollinate
